# ScriptsVBS/antigravity-cache-clean.service
[Unit]
Description=Clean Antigravity IDE cache on login
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'rm -rf ~/.config/Antigravity/Service\ Worker ~/.config/Antigravity/GPUCache ~/.config/Antigravity/Cache ~/.config/Antigravity/Code\ Cache 2>/dev/null'

[Install]
WantedBy=default.target
